# lab10/rpmbuild/BUILD/ssh-monitoring-1.0/ssh-monitoring.service
[Unit]
Description=ssh-monitoring
After=sshd.service


[Service]
Type=simple
ExecStart=/usr/bin/ssh-monitoring.sh
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target

